# rootfs/alpine/action-server
#!/sbin/openrc-run

description="starts the actix dmesg server"
pidfile="/dev/shm/${RC_SVCNAME}.pid"
supervisor="supervise-daemon"
command_background=true
environment
command="actionloop"

reload() {
  ebegin "Reloading ${RC_SVCNAME}"
  start-stop-daemon --signal HUP --pidfile "${pidfile}"
  eend $?
}
